# infrastructure/playbooks/roles/redeploy/templates/runner.service.j2
[Unit]
Description=IPFS benchmark runner
Requires=network.target
After=network.target
[Service]
ExecStart=/bin/bash -c "IPFS_PASSWORD=$(cat /etc/pwd.txt) /usr/bin/node {{ remote_folder }}/runner/index.js"
User={{ remote_user }}
Restart=always
Environment="STAGE=remote"
Environment="LOG_PRETTY=true"
Environment="BENCHMARK_USER={{ remote_user }}"
Environment="LOGLEVEL=debug"
Environment="CLINIC=true"
Environment="HOSTNAME=benchmarks.ipfs.team"
Environment="DATADIR=/data/tmp"
EnvironmentFile=-/etc/environment
[Install]
WantedBy=multi-user.target